# roles/satellite-audio/files/99-jabra.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{product}=="Jabra SPEAK 510 USB", RUN+="/usr/local/bin/jabra_connected_wrapper.sh"