# roles/hermes_docker/templates/hermes-docker-watchdog.service.j2
[Unit]
Description=Hermes Docker Tailscale watchdog
Requires=docker.service
After=docker.service tailscaled.service
ConditionPathExists={{ hermes_docker_watchdog_script_path }}

[Service]
Type=oneshot
ExecStart={{ hermes_docker_watchdog_script_path }}
